# systemd/hop-scan-all.service
[Unit]
Description=HoP repo scan (all repos)
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/head-of-product
EnvironmentFile=/opt/head-of-product/.env
ExecStart=/opt/head-of-product/.venv/bin/python -m agents.repo_analyst --all
TimeoutStartSec=3600
